# crawler/crontab
*/5 * * * * python3 /code/app/main.py >> /logfile
7 12 * * * python3 /code/app/main.py intraday >> /logfile
14 12 * * * python3 /code/app/main.py historic >> /logfile